# files/etc/init.d/set_up_huawei_modem
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1

start_service() {
	procd_open_instance "set_up_huawei_modem"
	procd_set_param command "/usr/sbin/set_up_huawei_modem"
	procd_close_instance
}
