# files/etc/hotplug.d/usb/01_set_up_huawei_modem
#!/bin/sh

# Find the Huawei modem by its vendor, product, and revision.
if [ "${PRODUCT}" = "12d1/1506/102" ] && [ "${ACTION}" = "add" ]; then
	set_up_huawei_modem
fi
